# modules/jira/templates/jira.erb
#!/bin/bash
#
# jira redhat init Nan Liu <nan@puppetlabs.com> 
# 
# chkconfig: - 60 50
# description: jira enterprise wiki (catalina tomcat)

# source function library
. /etc/rc.d/init.d/functions

# pull in sysconfig settings
[ -f /etc/sysconfig/jira ] && source /etc/sysconfig/jira

RETVAL=0
export JRE_HOME=/usr/lib/jvm/jre-1.6.0-openjdk
export CATALINA_PID='/var/run/jira.pid'
JIRA_HOME=<%= scope.lookupvar('params::jira_dir') -%>

start() {
  ${JIRA_HOME}/bin/startup.sh
  RETVAL=$?
}

stop() {
  ${JIRA_HOME}/bin/shutdown.sh
  RETVAL=$?
}

case "$1" in 
  start)
    start
  ;;
  stop)
    stop
  ;;
  restart)
    stop
    start
  ;;
  status)
    # this appears to be a bug with jira catalina.sh
    status -p ${CATALINA_PID}
    RETVAL=$?
  ;;
  *)
    echo "Usage: $(basename $0) {start|stop|status|restart}"
    exit 1
  ;;
esac

exit $RETVAL
